# Systemd/xlxd.timer
[Timer]
OnStartupSec=40

[Install]
WantedBy=multi-user.target
