# scripts/ambedwatchdog.service
[Unit]
Description=Check ambedwatchdog
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/ambedwatchdog
